# contrib/debian/infon-server.postinst
#!/bin/sh

#DEBHELPER#

/usr/sbin/update-metainit  >/dev/null
if [ -x "/etc/init.d/infon-server" ]; then

        if [ -x "`which invoke-rc.d 2>/dev/null`" ]; then
                invoke-rc.d infon-server start || exit $?
        else
		if [ -x /etc/init.d/infon-server ]; then
	                /etc/init.d/infon-server start || exit $?
		fi
        fi
fi

